# arquivo.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task" version="1.2">
    <Triggers>
        <CalendarTrigger>
            <Enabled>true</Enabled>
            <StartBoundary>2023-03-10T07:00:00-04:00</StartBoundary>
            <ScheduleByMonth>
                <DaysOfMonth>
                    <Day>1</Day>
                    <Day>2</Day>
                    <Day>3</Day>
                    <Day>4</Day>
                    <Day>5</Day>
                    <Day>6</Day>
                    <Day>7</Day>
                    <Day>8</Day>
                    <Day>9</Day>
                    <Day>10</Day>
                    <Day>11</Day>
                    <Day>12</Day>
                    <Day>13</Day>
                    <Day>14</Day>
                    <Day>15</Day>
                    <Day>16</Day>
                    <Day>17</Day>
                    <Day>18</Day>
                    <Day>19</Day>
                    <Day>20</Day>
                    <Day>21</Day>
                    <Day>22</Day>
                    <Day>23</Day>
                    <Day>24</Day>
                    <Day>25</Day>
                    <Day>26</Day>
                    <Day>27</Day>
                    <Day>28</Day>
                    <Day>29</Day>
                    <Day>30</Day>
                    <Day>31</Day>
                </DaysOfMonth>
                <Months>
                    <March/>
                    <April/>
                </Months>
            </ScheduleByMonth>
        </CalendarTrigger>
    </Triggers>
    <Actions Context="Author">
        <Exec>
            <Command>notepad.exe</Command>
        </Exec>
    </Actions>
</Task>