# arquivo.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task" version="1.2">
    <Triggers>
        <TimeTrigger>
            <Enabled>true</Enabled>
            <StartBoundary>2023-03-09T23:37:00-04:00</StartBoundary>
            <Repetition>
                <Interval>PT1M</Interval>
                <StopAtDurationEnd>false</StopAtDurationEnd>
            </Repetition>
        </TimeTrigger>
    </Triggers>
    <Actions Context="Author">
        <Exec>
            <Command>notepad.exe</Command>
        </Exec>
    </Actions>
</Task>