# arquivo.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task" version="1.2">
    <Triggers>
        <CalendarTrigger>
            <Repetition>
                <Interval>PT1M</Interval>
                <StopAtDurationEnd>true</StopAtDurationEnd>
                <Duration>P1D</Duration>
            </Repetition>
            <Enabled>true</Enabled>
            <StartBoundary>2023-03-10T17:23:11-04:00</StartBoundary>
            <ScheduleByMonthDayOfWeek>
                <DaysOfWeek>
                    <Monday/>
                    <Tuesday/>
                    <Wednesday/>
                    <Thursday/>
                    <Friday/>
                </DaysOfWeek>
                <Months>
                    <January/>
                </Months>
                <Weeks>
                    <Week>1</Week>
                    <Week>2</Week>
                    <Week>3</Week>
                    <Week>4</Week>
                </Weeks>
            </ScheduleByMonthDayOfWeek>
        </CalendarTrigger>
    </Triggers>
    <Actions Context="Author">
        <Exec>
            <Command>notepad.exe</Command>
        </Exec>
    </Actions>
</Task>