# arquivo.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task" version="1.2">
    <Triggers>
        <CalendarTrigger>
            <Enabled>true</Enabled>
            <StartBoundary>2023-03-10T16:00:00-04:00</StartBoundary>
            <ScheduleByMonth>
                <DaysOfMonth>
                    <Day>10</Day>
                </DaysOfMonth>
                <Months>
                    <February/>
                </Months>
            </ScheduleByMonth>
        </CalendarTrigger>
    </Triggers>
    <Actions Context="Author">
        <Exec>
            <Command>notepad.exe</Command>
        </Exec>
    </Actions>
</Task>