# arquivo.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task" version="1.2">
    <Triggers>
        <CalendarTrigger>
            <Repetition>
                <Interval>PT1M</Interval>
                <StopAtDurationEnd>true</StopAtDurationEnd>
                <Duration>P1D</Duration>
            </Repetition>
            <Enabled>true</Enabled>
            <StartBoundary>2023-03-10T15:25:31-04:00</StartBoundary>
            <ScheduleByMonth>
                <DaysOfMonth>
                    <Day>1</Day>
                    <Day>1</Day>
                    <Day>2</Day>
                    <Day>2</Day>
                    <Day>3</Day>
                    <Day>3</Day>
                    <Day>4</Day>
                    <Day>4</Day>
                    <Day>5</Day>
                    <Day>5</Day>
                </DaysOfMonth>
                <Months>
                    <January/>
                    <March/>
                </Months>
            </ScheduleByMonth>
        </CalendarTrigger>
    </Triggers>
    <Actions Context="Author">
        <Exec>
            <Command>notepad.exe</Command>
        </Exec>
    </Actions>
</Task>