# arquivo.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task" version="1.2">
    <Triggers>
        <CalendarTrigger>
            <Repetition>
                <Interval>PT1M</Interval>
                <Duration>PT1H</Duration>
                <StopAtDurationEnd>false</StopAtDurationEnd>
            </Repetition>
            <Enabled>true</Enabled>
            <StartBoundary>2023-03-10T01:00:00-04:00</StartBoundary>
            <ScheduleByMonth>
                <DaysOfMonth>
                    <Day>1</Day>
                    <Day>5</Day>
                </DaysOfMonth>
                <Months>
                    <January/>
                    <February/>
                    <March/>
                    <April/>
                    <May/>
                    <June/>
                    <July/>
                    <August/>
                    <September/>
                    <October/>
                    <November/>
                    <December/>
                </Months>
            </ScheduleByMonth>
        </CalendarTrigger>
        <CalendarTrigger>
            <Repetition>
                <Interval>PT1M</Interval>
                <Duration>PT1H</Duration>
                <StopAtDurationEnd>false</StopAtDurationEnd>
            </Repetition>
            <Enabled>true</Enabled>
            <StartBoundary>2023-03-10T01:00:00-04:00</StartBoundary>
            <ScheduleByMonth>
                <DaysOfMonth>
                    <Day>1</Day>
                    <Day>5</Day>
                </DaysOfMonth>
                <Months>
                    <January/>
                    <February/>
                    <March/>
                    <April/>
                    <May/>
                    <June/>
                    <July/>
                    <August/>
                    <September/>
                    <October/>
                    <November/>
                    <December/>
                </Months>
            </ScheduleByMonth>
        </CalendarTrigger>
        <CalendarTrigger>
            <Repetition>
                <Interval>PT1M</Interval>
                <Duration>PT1H</Duration>
                <StopAtDurationEnd>false</StopAtDurationEnd>
            </Repetition>
            <Enabled>true</Enabled>
            <StartBoundary>2023-03-10T02:00:00-04:00</StartBoundary>
            <ScheduleByMonth>
                <DaysOfMonth>
                    <Day>1</Day>
                    <Day>5</Day>
                </DaysOfMonth>
                <Months>
                    <January/>
                    <February/>
                    <March/>
                    <April/>
                    <May/>
                    <June/>
                    <July/>
                    <August/>
                    <September/>
                    <October/>
                    <November/>
                    <December/>
                </Months>
            </ScheduleByMonth>
        </CalendarTrigger>
        <CalendarTrigger>
            <Repetition>
                <Interval>PT1M</Interval>
                <Duration>PT1H</Duration>
                <StopAtDurationEnd>false</StopAtDurationEnd>
            </Repetition>
            <Enabled>true</Enabled>
            <StartBoundary>2023-03-10T02:00:00-04:00</StartBoundary>
            <ScheduleByMonth>
                <DaysOfMonth>
                    <Day>1</Day>
                    <Day>5</Day>
                </DaysOfMonth>
                <Months>
                    <January/>
                    <February/>
                    <March/>
                    <April/>
                    <May/>
                    <June/>
                    <July/>
                    <August/>
                    <September/>
                    <October/>
                    <November/>
                    <December/>
                </Months>
            </ScheduleByMonth>
        </CalendarTrigger>
        <CalendarTrigger>
            <Repetition>
                <Interval>PT1M</Interval>
                <Duration>PT1H</Duration>
                <StopAtDurationEnd>false</StopAtDurationEnd>
            </Repetition>
            <Enabled>true</Enabled>
            <StartBoundary>2023-03-10T03:00:00-04:00</StartBoundary>
            <ScheduleByMonth>
                <DaysOfMonth>
                    <Day>1</Day>
                    <Day>5</Day>
                </DaysOfMonth>
                <Months>
                    <January/>
                    <February/>
                    <March/>
                    <April/>
                    <May/>
                    <June/>
                    <July/>
                    <August/>
                    <September/>
                    <October/>
                    <November/>
                    <December/>
                </Months>
            </ScheduleByMonth>
        </CalendarTrigger>
        <CalendarTrigger>
            <Repetition>
                <Interval>PT1M</Interval>
                <Duration>PT1H</Duration>
                <StopAtDurationEnd>false</StopAtDurationEnd>
            </Repetition>
            <Enabled>true</Enabled>
            <StartBoundary>2023-03-10T03:00:00-04:00</StartBoundary>
            <ScheduleByMonth>
                <DaysOfMonth>
                    <Day>1</Day>
                    <Day>5</Day>
                </DaysOfMonth>
                <Months>
                    <January/>
                    <February/>
                    <March/>
                    <April/>
                    <May/>
                    <June/>
                    <July/>
                    <August/>
                    <September/>
                    <October/>
                    <November/>
                    <December/>
                </Months>
            </ScheduleByMonth>
        </CalendarTrigger>
        <CalendarTrigger>
            <Repetition>
                <Interval>PT1M</Interval>
                <Duration>PT1H</Duration>
                <StopAtDurationEnd>false</StopAtDurationEnd>
            </Repetition>
            <Enabled>true</Enabled>
            <StartBoundary>2023-03-10T04:00:00-04:00</StartBoundary>
            <ScheduleByMonth>
                <DaysOfMonth>
                    <Day>1</Day>
                    <Day>5</Day>
                </DaysOfMonth>
                <Months>
                    <January/>
                    <February/>
                    <March/>
                    <April/>
                    <May/>
                    <June/>
                    <July/>
                    <August/>
                    <September/>
                    <October/>
                    <November/>
                    <December/>
                </Months>
            </ScheduleByMonth>
        </CalendarTrigger>
        <CalendarTrigger>
            <Repetition>
                <Interval>PT1M</Interval>
                <Duration>PT1H</Duration>
                <StopAtDurationEnd>false</StopAtDurationEnd>
            </Repetition>
            <Enabled>true</Enabled>
            <StartBoundary>2023-03-10T04:00:00-04:00</StartBoundary>
            <ScheduleByMonth>
                <DaysOfMonth>
                    <Day>1</Day>
                    <Day>5</Day>
                </DaysOfMonth>
                <Months>
                    <January/>
                    <February/>
                    <March/>
                    <April/>
                    <May/>
                    <June/>
                    <July/>
                    <August/>
                    <September/>
                    <October/>
                    <November/>
                    <December/>
                </Months>
            </ScheduleByMonth>
        </CalendarTrigger>
        <CalendarTrigger>
            <Repetition>
                <Interval>PT1M</Interval>
                <Duration>PT1H</Duration>
                <StopAtDurationEnd>false</StopAtDurationEnd>
            </Repetition>
            <Enabled>true</Enabled>
            <StartBoundary>2023-03-10T05:00:00-04:00</StartBoundary>
            <ScheduleByMonth>
                <DaysOfMonth>
                    <Day>1</Day>
                    <Day>5</Day>
                </DaysOfMonth>
                <Months>
                    <January/>
                    <February/>
                    <March/>
                    <April/>
                    <May/>
                    <June/>
                    <July/>
                    <August/>
                    <September/>
                    <October/>
                    <November/>
                    <December/>
                </Months>
            </ScheduleByMonth>
        </CalendarTrigger>
        <CalendarTrigger>
            <Repetition>
                <Interval>PT1M</Interval>
                <Duration>PT1H</Duration>
                <StopAtDurationEnd>false</StopAtDurationEnd>
            </Repetition>
            <Enabled>true</Enabled>
            <StartBoundary>2023-03-10T05:00:00-04:00</StartBoundary>
            <ScheduleByMonth>
                <DaysOfMonth>
                    <Day>1</Day>
                    <Day>5</Day>
                </DaysOfMonth>
                <Months>
                    <January/>
                    <February/>
                    <March/>
                    <April/>
                    <May/>
                    <June/>
                    <July/>
                    <August/>
                    <September/>
                    <October/>
                    <November/>
                    <December/>
                </Months>
            </ScheduleByMonth>
        </CalendarTrigger>
    </Triggers>
    <Actions Context="Author">
        <Exec>
            <Command>notepad.exe</Command>
        </Exec>
    </Actions>
</Task>